# deploy/launchd/com.worxphere.notion-api-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.worxphere.notion-api-gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>notion-gateway</string>
    <string>poll</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/agent/project/notion-api-gateway-py</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>/Users/agent/project/notion-api-gateway-py/operations/logs/poll.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/agent/project/notion-api-gateway-py/operations/logs/poll.err.log</string>
</dict>
</plist>
