# deploy/notion-gateway.service
[Unit]
Description=Notion API Gateway Poller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/app-notion/notion-api-gateway-py
Environment=PATH=/home/app-notion/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/app-notion
Environment=AWS_REGION=ap-northeast-2
ExecStart=/home/app-notion/bin/uv run notion-gateway poll
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=notion-gateway

[Install]
WantedBy=default.target
